# scripts/config/celery-beat.service
[Unit]
Description=Celery Beat for SWATGenX
After=network.target redis-server.service
Requires=redis-server.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=${WEBAPP_DIR}
Environment="PATH=${BASE_DIR}/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
Environment="PYTHONPATH=${BASE_DIR}:${WEBAPP_DIR}"
EnvironmentFile=-${CONFIG_DIR}/swatgenx.env
ExecStart=${BASE_DIR}/.venv/bin/celery -A celery_worker beat --loglevel=info

Restart=always
RestartSec=5
StandardOutput=append:${LOG_DIR}/celery-beat.log
StandardError=append:${LOG_DIR}/celery-beat-error.log

[Install]
WantedBy=multi-user.target
