# scripts/config/flask-app.service
[Unit]
Description=SWATGenX Flask Web Application
After=network.target redis-server.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=${WEBAPP_DIR}
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:${BASE_DIR}/.venv/bin"
Environment="FLASK_APP=run.py"
Environment="FLASK_ENV=production"
Environment="PYTHONPATH=${BASE_DIR}:${WEBAPP_DIR}"
Environment="REDIS_URL=redis://127.0.0.1:6379/0"
ExecStart=${BASE_DIR}/.venv/bin/gunicorn --workers 4 --bind 127.0.0.1:5050 run:app --timeout 120 --access-logfile ${LOG_DIR}/gunicorn-access.log --error-logfile ${LOG_DIR}/gunicorn-error.log

# Restart policy
Restart=always
RestartSec=5

# Logging
StandardOutput=append:${LOG_DIR}/flask-app.log
StandardError=append:${LOG_DIR}/flask-app-error.log

[Install]
WantedBy=multi-user.target
