# scripts/celery-services/services/redis-fix.service
[Unit]
Description=Redis WRONGTYPE Error Repair Service
After=redis.service
Requires=redis.service

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/data/SWATGenXApp/codes/scripts
ExecStart=/data/SWATGenXApp/codes/scripts/fix_redis_wrongtype_batch_wrapper.sh --yes --batch-size 50 --timeout 7200

# Logging
StandardOutput=append:/data/SWATGenXApp/codes/web_application/logs/redis-fix-service.log
StandardError=append:/data/SWATGenXApp/codes/web_application/logs/redis-fix-service-error.log

# Safety limits
TimeoutStartSec=7200
TimeoutStopSec=300
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
