# scripts/celery-worker.service
[Unit]
Description=Celery Worker for Flask Application
After=network.target redis.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/data/SWATGenXApp/codes/web_application
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/data/SWATGenXApp/codes/.venv/bin"
Environment="FLASK_ENV=production"
Environment="PYTHONPATH=/data/SWATGenXApp/codes:/data/SWATGenXApp/codes/web_application"
Environment="C_FORCE_ROOT=true"
Environment="CELERY_LOG_LEVEL=DEBUG"

# Use celery_app instead of celery_worker for the Celery instance
ExecStart=/data/SWATGenXApp/codes/.venv/bin/celery -A celery_app worker --loglevel=debug --concurrency=4

# Restart policy
Restart=always
RestartSec=5

# Increase timeout for service startup
TimeoutStartSec=30

# Logging
StandardOutput=append:/data/SWATGenXApp/codes/web_application/logs/celery-worker.log
StandardError=append:/data/SWATGenXApp/codes/web_application/logs/celery-worker-error.log

[Install]
WantedBy=multi-user.target