# scripts/celery-worker.service
[Unit]
Description=Celery Worker for SWATGenX
After=network.target redis-server.service
Requires=redis-server.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/data/SWATGenXApp/codes/web_application
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/data/SWATGenXApp/codes/.venv/bin"
Environment="FLASK_ENV=production"
Environment="PYTHONPATH=/data/SWATGenXApp/codes:/data/SWATGenXApp/codes/web_application"
Environment="C_FORCE_ROOT=true"
Environment="CELERY_LOG_LEVEL=DEBUG"
Environment="REDIS_URL=redis://127.0.0.1:6379/0"
Environment="API_SERVER_PORT=5050"

# Configure worker with fixed parameters for stable queue processing
ExecStart=/data/SWATGenXApp/codes/.venv/bin/celery -A celery_worker worker --loglevel=info -n model_worker@%%h -Q model_creation --concurrency=4 --prefetch-multiplier=1 --max-tasks-per-child=10 -O fair --without-heartbeat --without-gossip

# Proper restart policy to handle connection issues
Restart=always
RestartSec=10

# Set resource limits
LimitNOFILE=65536
TimeoutStartSec=60
TimeoutStopSec=60

# Logging
StandardOutput=append:/data/SWATGenXApp/codes/web_application/logs/celery-worker.log
StandardError=append:/data/SWATGenXApp/codes/web_application/logs/celery-worker-error.log

[Install]
WantedBy=multi-user.target