# scripts/celery-multi-worker.service
[Unit]
Description=Celery Multiple Workers Service
After=network.target redis.service

[Service]
Type=forking
User=www-data
Group=www-data
WorkingDirectory=/data/SWATGenXApp/codes/web_application
Environment=FLASK_APP=run.py
Environment=PYTHONPATH=/data/SWATGenXApp/codes:/data/SWATGenXApp/codes/SWATGenX

# Set Celery configuration via environment variables
Environment=CELERY_WORKER_COUNT=30
Environment=CELERY_WORKER_CONCURRENCY=8
Environment=CELERY_MAX_TASKS_PER_CHILD=100
Environment=CELERY_MAX_MEMORY_PER_CHILD_MB=8192
Environment=CELERY_WORKER_PREFETCH_MULTIPLIER=8
Environment=CELERY_TASK_SOFT_TIME_LIMIT=43200
Environment=CELERY_TASK_TIME_LIMIT=86400
Environment=CELERY_DISABLE_RATE_LIMITS=true
Environment=CELERY_TASK_DEFAULT_RATE_LIMIT=
Environment=CELERY_BROKER_CONNECTION_RETRY=true
Environment=CELERY_BROKER_CONNECTION_MAX_RETRIES=20
Environment=CELERY_REDIS_MAX_CONNECTIONS=500
Environment=CELERY_MODEL_CREATION_WORKER_PERCENT=70

ExecStart=/data/SWATGenXApp/codes/scripts/restart_celery_workers.sh
Restart=always
RestartSec=30s
TimeoutStartSec=180
TimeoutStopSec=180
LimitNOFILE=65536

# Resource limits - remove CPU limits
LimitCPU=infinity
LimitRSS=infinity
LimitFSIZE=infinity
LimitNPROC=32768

[Install]
WantedBy=multi-user.target