# scripts/celery-multi-worker.service
[Unit]
Description=Celery Multiple Workers Service
After=network.target redis.service

[Service]
Type=forking
User=www-data
Group=www-data
WorkingDirectory=/data/SWATGenXApp/codes/web_application
Environment=FLASK_APP=run.py
Environment=PYTHONPATH=/data/SWATGenXApp/codes:/data/SWATGenXApp/codes/SWATGenX
# Set MAX_WORKER_COUNT to 0 to auto-calculate based on CPU count
Environment=MAX_WORKER_COUNT=0
Environment=WORKER_CONCURRENCY=0
ExecStart=/data/SWATGenXApp/codes/scripts/restart_celery_workers.sh
Restart=always
RestartSec=30s
TimeoutStartSec=180
TimeoutStopSec=180
LimitNOFILE=65536

# Resource limits - remove CPU limits
LimitCPU=infinity
LimitRSS=infinity
LimitFSIZE=infinity
LimitNPROC=32768

[Install]
WantedBy=multi-user.target